# vm_profiler.service
[Unit]
Description=VM Profiler Service

[Service]
ExecStart=~/Desktop/git/vm_profiler/runAllTests.bash

[Install]
WantedBy=multi-user.target